# gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn.sock
User=root

[Install]
WantedBy=sockets.target
